# .config/autostart/kbd-capslk.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=kbd-capslk
Exec=kbd-capslk.sh
Terminal=false
Categories=kbd
